# deploy/systemd/booking-platform.service.example
[Unit]
Description=booking-platform service
After=network.target

[Service]
Type=simple
User=booking
Group=booking
WorkingDirectory=/opt/booking-platform/current
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
